# Task_8_Crontab/crontab
SHELL=/bin/bash
MAILTO=ragen
*/2 * * * * ~/Programs/BashProgr/TimeAndDate.sh >> /tmp/cron_test

